# ansible/roles/llm_api/templates/llm_api.service.j2
[Unit]
Description=LLM API
After=network.target

[Service]
User={{ llm_api_user | default("ec2-user") }}
Group={{ llm_api_group | default("ec2-user") }}
WorkingDirectory={{ llm_api_workdir }}
Environment="GRACEFUL_TIMEOUT={{ llm_api_graceful_timeout | default(25) }}"
ExecStart={{ llm_api_exec }}  # 例: /usr/local/bin/uvicorn llm_api.app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=30
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
